# cron.d/virsh-wrapper
*/2 * * * * wrapper test -x /usr/local/bin/virsh_wrapper.sh && sudo /usr/local/bin/virsh_wrapper.sh
