# backend/cron.d/virsh-wrapper
*/2 * * * * wrapper test -x /usr/local/bin/virsh_wrapper/backend/virsh_wrapper.sh && /usr/local/bin/virsh_wrapper/backend/virsh_wrapper.sh
